# deploy/quadlet/void-frontend.container
[Unit]
Description=Void Frontend
Requires=void-network.service void-backend.service
After=void-network.service void-backend.service

[Container]
ContainerName=frontend
Image=ghcr.io/scetrov/void-eid/void-eid-frontend:main
Pull=newer
PublishPort=5173:8080
AutoUpdate=registry
Network=void.network
GroupAdd=keep-groups
EnvironmentFile=%h/.config/void-eid/.env

[Service]
Restart=always
StartLimitBurst=10
RestartSec=30

[Install]
WantedBy=default.target
